# deploy/systemd/dubai-garments-storefront.service
[Unit]
Description=Dubai Garments Storefront (Next.js)
After=network.target dubai-garments-fastapi.service
Wants=dubai-garments-fastapi.service

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__PROJECT_ROOT__/apps/storefront-dubai_garments
Environment=NODE_ENV=production
EnvironmentFile=__PROJECT_ROOT__/apps/storefront-dubai_garments/.env.local
ExecStart=/usr/bin/env npm run start -- -H 0.0.0.0 -p 3000
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
